# apk_repo_upgrade.timer
[Unit]
Description=Regular apk repo upgrade

[Timer]
OnUnitActiveSec=6d
OnBootSec=10min

[Install]
WantedBy=timers.target
