# phone_upgrade.timer
[Unit]
Description=Regular phone pkg upgrade 

[Timer]
OnUnitActiveSec=2w 
OnBootSec=10min

[Install]
WantedBy=timers.target
